# .config/autostart/home-update.desktop
[Desktop Entry]
Name=Home Backup
Comment=Does a backup each startup.
Type=Application
Terminal=false
Exec=/home/jsd/.local/bin/home-update.sh